# x.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/xinit -- -nocursor
Restart=on-failure

[Install]
WantedBy=multi-user.target
