# x.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/xinit -- -nocursor
Restart=always
RestartSec=120

[Install]
WantedBy=default.target
